# deploy/systemd/oqim-api.service
[Unit]
Description=OQIM Backend API
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=oqim
Group=oqim
WorkingDirectory=/opt/oqim/current/backend
EnvironmentFile=/etc/oqim/oqim.env
ExecStart=/opt/oqim/current/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/oqim/current /var/lib/oqim /var/log/oqim

[Install]
WantedBy=multi-user.target
